# s6/tailscaled/run
#!/command/with-contenv sh
# tailscaled — uses real kernel TUN when the host allows it (Fly.io VMs), and
# falls back to userspace networking when it doesn't (Railway/most PaaS).
#
# Kernel TUN matters: `tailscale serve`'s userspace (gVisor) netstack stalls TLS
# handshakes (slow + wedges after a couple connections). With a real TUN device
# inbound traffic is handled by the kernel, so connections are fast and stable.
#
#   --statedir → on-disk "var root", required for `tailscale serve` TLS certs
#                (mem: state => "no TailscaleVarRoot"). Mount a volume here to
#                keep node identity + cert across restarts.
# /run is often a runtime tmpfs, so (re)create the socket dir before launching.
mkdir -p /var/run/tailscale /var/lib/tailscale /dev/net

# Create the TUN device node if the kernel module is present but /dev isn't
# populated (common in minimal VM images).
if [ ! -c /dev/net/tun ]; then
    mknod /dev/net/tun c 10 200 2>/dev/null || true
    chmod 600 /dev/net/tun 2>/dev/null || true
fi

if [ -c /dev/net/tun ]; then
    echo "[tailscaled] /dev/net/tun present — kernel TUN mode" >&2
    set -- --statedir=/var/lib/tailscale --socket=/var/run/tailscale/tailscaled.sock
else
    echo "[tailscaled] no TUN device — userspace networking mode" >&2
    set -- --tun=userspace-networking --statedir=/var/lib/tailscale --socket=/var/run/tailscale/tailscaled.sock
fi

exec /usr/local/bin/tailscaled "$@"
